# whosatmyfeeder/rootfs/etc/services.d/mqttclient/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the mqttclient service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
bashio::log.info "Starting mqttclient..."

## Convert the config file from JSON to YAML
bashio::log.info "Converting config from JSON to YAML (FIXME: runs twice)..."
python3 -c 'import sys, yaml, json; yaml.safe_dump(json.load(sys.stdin), sys.stdout, default_flow_style=False)' < /data/options.json > /data/config.yaml

## Run your program
bashio::log.info "Using config file at ${CONFIG_PATH}"
exec python3 -u /usr/app/speciesid.py
