# matrix-ng/rootfs/etc/s6-overlay/s6-rc.d/client-element/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start client-element service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
bashio::log.info "Starting client-element..."

export DEFAULT_HS_URL=$(bashio::config 'hs_url')
export DEFAULT_IS_URL=$(bashio::config 'is_url')

## Run your program
element server
